# fs2/usr/lib/ipkg/info/dph151-misc.list
/opt/oem/diagfiles
/opt/oem/basic_oam.dat
/etc/rcS.d/S06checkipconfig
/etc/init.d/checkipconfig
/opt/oem/ipconfig.dat
