# fs2/usr/lib/ipkg/info/etc.list
/opt/ipaccess/bin/getbank
/etc/fstab
/etc/init.d/mountoemparts
/etc/init.d/ubootenvcheck
/etc/rcS.d/S03ubootenvcheck
