# fs2/usr/lib/ipkg/info/initscripts.list
/opt/ipaccess/bin/setnv_env.sh
/etc/init.d/cpu_alignment
/etc/init.d/dhcp
/etc/init.d/rc.sysinit
/opt/ipaccess/Utils/scripts/cdrlogcopy
/etc/init.d/picostatus
/opt/ipaccess/ipconfig.dat
/opt/ipaccess/Drivers/cpetd_load
/etc/init.d/functions
/etc/init.d/start
/etc/init.d/loadmodules
/opt/ipaccess/Utils/scripts/save_reboot_cause
/opt/ipaccess/managelogs_default.cfg
/opt/ipaccess/Utils/scripts/set_reboot_cause_oem
/etc/init.d/syslog
/etc/init.d/rc.local
/etc/init.d/tcpdump
/etc/init.d/dmistart
/etc/init.d/l1apiApp
/opt/ipaccess/Utils/scripts/apdiag.dmi
/etc/init.d/random
/opt/ipaccess/bin/managelogs
/opt/ipaccess/Utils/scripts/set_reboot_cause
/etc/init.d/rcS
/etc/init.d/opnwl
/etc/init.d/last_reboot
/opt/ipaccess/Utils/scripts/old_uboot.env
/opt/ipaccess/bin/readtemp
/etc/init.d/fallbackdate
/opt/ipaccess/Utils/scripts/factoryrestore
/opt/ipaccess/Utils/scripts/getVarVal
/opt/ipaccess/Utils/scripts/reboot
/etc/init.d/startgetty
/opt/ipaccess/Utils/scripts/setVarVal
/opt/ipaccess/Utils/scripts/relocate_uboot_undo
/etc/init.d/opcne
/etc/init.d/mountoemparts
/opt/ipaccess/Utils/scripts/setup_vlan.sh
/etc/init.d/opnormal
/etc/profile
/opt/ipaccess/cron/crontabs/root
/opt/ipaccess/bin/cramfsmount
/opt/ipaccess/bin/progpartition
/opt/ipaccess/bin/updateresolv
/opt/ipaccess/Drivers/cpetd_unload
/opt/ipaccess/Utils/scripts/relocate_uboot
/etc/init.d/crond
/opt/ipaccess/bin/get_reverse_ssh_key
/etc/init.d/opdtcs
/opt/ipaccess/bin/purgepmlogs
/etc/init.d/verifyimages
/etc/init.d/startmode
/opt/ipaccess/bin/waitforpids
/etc/init.d/ubootenvcheck
/etc/init.d/picoinit
/etc/init.d/nwlIf
/opt/ipaccess/Utils/scripts/reboot_cause_codes
/opt/ipaccess/Utils/scripts/ntplogger
/opt/ipaccess/bin/reverse_ssh
/opt/ipaccess/bin/checkntp.sh
/opt/ipaccess/Utils/scripts/fsvalidator.sh
/etc/init.d/end
/opt/ipaccess/Utils/scripts/post_swdl_hook
/opt/ipaccess/Utils/scripts/servicecounter
/opt/ipaccess/bin/diagupload
/opt/ipaccess/Utils/scripts/getbanner
/etc/inittab
