# fs2/usr/lib/ipkg/info/iptables.list
/usr/sbin/iptables-restore
/etc/sysconfig/iptables-205J-rules
/etc/sysconfig/iptables-205I-rules
/usr/sbin/iptables
/etc/rcS.d/S06iptables
/etc/sysconfig/iptables-224F-rules
/etc/sysconfig/iptables-224E-rules
/etc/sysconfig/iptables-nano-rules
/usr/sbin/iptables-static
/opt/ipaccess/bin/fwcontrol
/usr/share/doc/iptables-1.4.2/COPYING
/etc/sysconfig/iptables-224J-rules
/etc/sysconfig/iptables-224I-rules
/etc/sysconfig/iptables-customerA-rules
/etc/sysconfig/iptables-205F-rules
/etc/sysconfig/iptables-205E-rules
/etc/init.d/iptablesinit
/usr/sbin/iptables-save
/etc/sysconfig/iptables-rules
