# fs2/usr/lib/ipkg/info/procmon-pc302.list
/etc/rcS.d/S85startprocmon
/opt/ipaccess/bin/procmon
/opt/ipaccess/bin/procmon.579.11.144
/opt/ipaccess/procmon.conf
/etc/init.d/startprocmon
