# fs2/usr/lib/ipkg/info/procmon-pc302.list
/opt/ipaccess/bin/procmon.563.21.2
/etc/rcS.d/S85startprocmon
/opt/ipaccess/bin/procmon
/opt/ipaccess/procmon.conf
/etc/init.d/startprocmon
