# fs2/etc/init.d/rc.local
#!/bin/bash
#
# add local initialisation

