# fs1/cramfs-root/usr/lib/ipkg/info/dph151-leds.list
/opt/cisco/dph151/ledsequences
/etc/init.d/loadleds
/etc/rcS.d/S35loadleds
