# fs1/usr/lib/ipkg/info/cisco-common.list
/opt/oem/diagfiles
/opt/cisco/leds/ledsequences-281
/opt/cisco/leds/ledsequences
/etc/init.d/loadleds
/etc/rcS.d/S35loadleds
