# fs2/opt/ipaccess/cron/crontabs/root
0 * * * * /opt/ipaccess/Utils/scripts/servicecounter
1 0 * * * /opt/ipaccess/Utils/scripts/ntplogger
0 */3 * * * /opt/ipaccess/Utils/scripts/cdrlogcopy periodic
